# scripts/zeroclaw.service
[Unit]
Description=ZeroClaw Autonomous Agent Runtime
Documentation=https://github.com/zeroclaw/zeroclaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ZEROCLAW_USER
WorkingDirectory=ZEROCLAW_HOME/code/zeroclaw
ExecStart=ZEROCLAW_HOME/code/zeroclaw/target/release/zeroclaw daemon
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
Environment=PATH=ZEROCLAW_HOME/.local/bin:ZEROCLAW_HOME/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=ZEROCLAW_HOME

[Install]
WantedBy=multi-user.target
